# scripts/chat2api.service
[Unit]
Description=Chat2API Service
After=network.target

[Service]
Type=simple
User=admin
Group=Users
WorkingDirectory=/opt/Chat2API
Environment=DISPLAY=:99
Environment=ELECTRON_IS_DEV=0
Environment=NODE_ENV=production
Environment=HOME=/home/admin
ExecStart=/usr/local/bin/chat2api-launcher
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Allow X11 access
Environment=XAUTHORITY=/home/admin/.Xauthority

[Install]
WantedBy=multi-user.target
